# rootfs/etc/services.d/phone-system/run
#!/usr/bin/with-contenv bashio

bashio::log.info "================================================"
bashio::log.info "  Advanced Phone System - Starting..."
bashio::log.info "================================================"

# Load configuration
SIP_ENABLED=$(bashio::config 'sip_trunk.enabled')
SIP_PROVIDER=$(bashio::config 'sip_trunk.provider')
SIP_HOST=$(bashio::config 'sip_trunk.host')
SIP_PORT=$(bashio::config 'sip_trunk.port')
SIP_USER=$(bashio::config 'sip_trunk.username')
SIP_PASS=$(bashio::config 'sip_trunk.password')
SIP_DOMAIN=$(bashio::config 'sip_trunk.from_domain')
SIP_CALLER=$(bashio::config 'sip_trunk.caller_number')

# Get recording system configuration
RECORDING_EXTENSION=$(bashio::config 'recording_extension')
RECORDING_PASSWORD=$(bashio::config 'recording_password')

# Create directories
mkdir -p /var/lib/asterisk/sounds/custom
mkdir -p /var/spool/asterisk/outgoing
mkdir -p /data/recordings
mkdir -p /data/database

# Generate Asterisk sip.conf
bashio::log.info "Generating Asterisk configuration..."

cat > /etc/asterisk/sip.conf <<'EOF'
[general]
context=default
allowoverlap=no
udpbindaddr=0.0.0.0
tcpenable=yes
tcpbindaddr=0.0.0.0
transport=udp
srvlookup=yes
allowguest=no
alwaysauthreject=yes
nat=force_rport,comedia
externrefresh=10

; Codecs
disallow=all
allow=ulaw
allow=alaw
allow=g729
allow=gsm

; RTP Settings
rtpstart=10000
rtpend=10099

; Security
requirecalltoken=yes

EOF

# Add SIP trunk if enabled
if [ "$SIP_ENABLED" = "true" ]; then
    bashio::log.info "Configuring SIP trunk: ${SIP_PROVIDER}"
    
    # Use from_domain if provided, otherwise use host
    DOMAIN="${SIP_DOMAIN:-$SIP_HOST}"
    
    # Get caller number if provided
    SIP_CALLER=$(bashio::config 'sip_trunk.caller_number')
    
    # Build registration string separately to avoid escape issues
    REG_STRING="${SIP_USER}:${SIP_PASS}@${SIP_HOST}:${SIP_PORT}/${SIP_USER}"
    
    # Write SIP trunk config - removed fromuser line
    {
        echo ""
        echo "; SIP Trunk Registration"
        echo "register => ${REG_STRING}"
        echo ""
        echo "[trunk_main]"
        echo "type=peer"
        echo "host=${SIP_HOST}"
        echo "port=${SIP_PORT}"
        echo "defaultuser=${SIP_USER}"
        echo "secret=${SIP_PASS}"
        echo "fromdomain=${DOMAIN}"
        echo "insecure=port,invite"
        echo "context=inbound"
        echo "dtmfmode=rfc2833"
        echo "canreinvite=no"
        echo "qualify=yes"
        echo "nat=force_rport,comedia"
        echo ""
    } >> /etc/asterisk/sip.conf
    bashio::log.info "✓ SIP trunk configured (user: ${SIP_USER})"
    bashio::log.info "  Registration string: ${REG_STRING}"
else
    bashio::log.warning "SIP trunk disabled - enable in configuration to make calls"
fi

# Add extensions
bashio::log.info "Configuring extensions..."
EXTENSION_COUNT=$(bashio::config 'extensions|length')

for (( i=0; i < EXTENSION_COUNT; i++ )); do
    EXT_NUMBER=$(bashio::config "extensions[${i}].number")
    EXT_NAME=$(bashio::config "extensions[${i}].name")
    EXT_SECRET=$(bashio::config "extensions[${i}].secret")
    
    bashio::log.info "  - Extension ${EXT_NUMBER}: ${EXT_NAME}"
    
    cat >> /etc/asterisk/sip.conf <<EOF

[${EXT_NUMBER}]
type=friend
secret=${EXT_SECRET}
context=internal
host=dynamic
dtmfmode=rfc2833
canreinvite=no
nat=force_rport,comedia
qualify=yes
callgroup=1
pickupgroup=1
callerid="${EXT_NAME}" <${EXT_NUMBER}>

EOF
done

bashio::log.info "✓ ${EXTENSION_COUNT} extension(s) configured"
bashio::log.info "Recording extension configured: ${RECORDING_EXTENSION}"

# Add recording extension context to extensions.conf
bashio::log.info "Configuring recording system..."
cat >> /etc/asterisk/extensions.conf << EOF

; ============================================
; Recording System
; ============================================
[record-messages]
exten => s,1,NoOp(=== Recording System ===)
 same => n,Answer()
 same => n,Wait(1)
 same => n,Playback(hello)
 same => n,Playback(please-enter-the)
 same => n,Read(PASSWORD,password,4)
 same => n,GotoIf(\$[\${PASSWORD} = ${RECORDING_PASSWORD}]?authorized:unauthorized)
 same => n(unauthorized),Playback(invalid)
 same => n,Hangup()
 same => n(authorized),Playback(accepted)
 same => n(get_name),Playback(please-enter-your&extension)
 same => n,Read(RECORDING_ID,extension,10)
 same => n,GotoIf(\$["\${RECORDING_ID}" = ""]?get_name)
 same => n,Playback(at-tone-rec-msg)
 same => n,Playback(beep)
 same => n,Set(TIMESTAMP=\${STRFTIME(\${EPOCH},,%Y%m%d_%H%M%S)})
 same => n,Set(FILENAME=custom_\${RECORDING_ID}_\${TIMESTAMP})
 same => n,Record(/var/lib/asterisk/sounds/custom/\${FILENAME}.wav,3,30,y)
 same => n,System(curl -X POST "http://localhost:8088/api/recordings/register?filename=\${FILENAME}.wav&recording_id=\${RECORDING_ID}" 2>/dev/null)
 same => n,Playback(thank-you-for-calling)
 same => n,Hangup()

EOF

bashio::log.info "✓ Recording system configured"

# Set permissions
chown -R asterisk:asterisk /etc/asterisk
chown -R asterisk:asterisk /var/lib/asterisk
chown -R asterisk:asterisk /var/spool/asterisk

# Kill any existing Asterisk processes
bashio::log.info "Cleaning up any existing Asterisk processes..."
pkill -9 asterisk 2>/dev/null || true
rm -f /var/run/asterisk/asterisk.ctl 2>/dev/null || true
sleep 2

# Start Asterisk
bashio::log.info "Starting Asterisk PBX..."
asterisk -f &
ASTERISK_PID=$!

# Wait for Asterisk to initialize
sleep 5

# Check Asterisk status (BusyBox compatible)
if kill -0 $ASTERISK_PID 2>/dev/null; then
    bashio::log.info "✓ Asterisk PBX started successfully"
else
    bashio::log.error "✗ Asterisk failed to start!"
    exit 1
fi

# Start API service
bashio::log.info "Starting API service..."
cd /app
python3 api_service.py &
API_PID=$!

# Wait for API to start
sleep 3

# Check API status (BusyBox compatible)
if kill -0 $API_PID 2>/dev/null; then
    bashio::log.info "✓ API service started successfully"
else
    bashio::log.error "✗ API service failed to start!"
    kill $ASTERISK_PID 2>/dev/null
    exit 1
fi

bashio::log.info "================================================"
bashio::log.info "  Advanced Phone System is READY!"
bashio::log.info "  Web UI: http://homeassistant.local:8088"
bashio::log.info "================================================"

# Monitor both processes (BusyBox compatible)
while true; do
    if ! kill -0 $ASTERISK_PID 2>/dev/null; then
        bashio::log.error "Asterisk process died - restarting addon..."
        kill $API_PID 2>/dev/null
        exit 1
    fi
    
    if ! kill -0 $API_PID 2>/dev/null; then
        bashio::log.error "API service died - restarting addon..."
        kill $ASTERISK_PID 2>/dev/null
        exit 1
    fi
    
    sleep 10
done
